# support/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
user=root
logfile_maxbytes = 0
pidfile=/supervisord.pid

[program:rebuilder]
command=/usr/local/bin/watchandrebuild.sh /data /var/www/html /usr/local/bin/make-rfc-from-md.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;" -c /etc/nginx/nginx.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

